# crontab
* * * * * cd /home/upp/upp && php artisan schedule:run >> /dev/null 2>&1
